# elastic-network-interfaces.conf
# This task finds and configures elastic network interfaces
# left in an unconfigured state.

start on started network
task
exec /sbin/ec2ifscan
